# vagrant/hack/ansible/roles/kubernetes-vagrant/templates/etc/systemd/system/kubelet.service.d/09-extra-args.conf
[Service]
Environment="KUBELET_EXTRA_ARGS=--node-ip={{ node_ip }}{% for k, v in kubeletConfig_extra_args.items() %} --{{k}}='{{v}}'{%- endfor %}"
